# deploy/evening-bot.service
[Unit]
Description=Evening Learning Telegram Bot
After=network.target evening-api.service
Wants=evening-api.service

[Service]
User=deploy
Group=deploy
WorkingDirectory=/opt/evening-learning/backend
Environment="PATH=/opt/evening-learning/backend/venv/bin"
EnvironmentFile=/opt/evening-learning/backend/.env
ExecStart=/opt/evening-learning/backend/venv/bin/python -m app.bot_polling
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
